# docker/gateway/etc/service/postgresql/finish
#!/bin/bash
echo "Stopping PostgreSQL..."
su postgres -c "pg_ctl -D /var/lib/postgresql/data stop -m fast"
